# biglinux-livecd/usr/lib/systemd/system/biglinux-language-suggestion.service
[Unit]
Description=Detect the likely BigLinux live-session language
Documentation=https://github.com/biglinux/biglinux-livecd
ConditionPathExists=/livefs-pkgs.txt
After=local-fs.target NetworkManager.service systemd-networkd.service
Before=display-manager.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/python /usr/lib/biglinux-livecd/language_suggestion_probe.py
TimeoutStartSec=3s
RemainAfterExit=true
RuntimeDirectory=biglinux-language-probe
RuntimeDirectoryMode=0755
UMask=0077
PrivateMounts=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/biglinux-language-probe
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true
Environment=PYTHONDONTWRITEBYTECODE=1

[Install]
WantedBy=graphical.target
