# biglinux-livecd/usr/lib/systemd/system/livecd-tweaks.service
# BigLinux LiveCD Tweaks Service
# Applies system optimizations during LiveCD boot for better performance

[Unit]
Description=BigLinux LiveCD System Tweaks and Optimizations
Documentation=https://github.com/biglinux
Before=display-manager.service mhwd-live.service
ConditionPathExists=/livefs-pkgs.txt

[Service]
Type=oneshot
ExecStart=/usr/bin/livecd-tweaks
RemainAfterExit=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
