# biglinux-livecd/usr/lib/systemd/system/livecd-tweaks.service
# BigLinux LiveCD Tweaks Service
# Applies system optimizations during LiveCD boot for better performance

[Unit]
Description=BigLinux LiveCD System Tweaks and Optimizations
Documentation=https://github.com/biglinux
Before=mhwd-live.service

[Service]
Type=oneshot
ExecStart=/usr/bin/livecd-tweaks
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
