# pkgbuild/biglinux-livecd.install
#!/usr/bin/env bash

post_install() {
    # Enable livecd-tweaks service for boot optimizations
    systemctl enable livecd-tweaks.service
}
