# pkgbuild/biglinux-livecd.install
post_install() { 
    systemctl enable livecd-tweaks.service
}
