# pkgbuild/biglinux-livecd.install
post_install() { 
systemctl enable mhwd-fix.service
}
