# biglinux-livecd/etc/xdg/autostart/live-big-autostart.desktop
[Desktop Entry]
Name=live-big-autostart
Exec=/usr/bin/live-big-autostart
Terminal=false
Type=Application
Categories=
